# GOrbitS/ci/gorbits.service.example
# Ejemplo systemd para el servidor GOrbitS (copiar a /etc/systemd/system/gorbits.service).
# Ajustar User, WorkingDirectory y rutas según vuestro servidor.
#
#   sudo cp gorbits.service.example /etc/systemd/system/gorbits.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable gorbits
#   sudo systemctl start gorbits

[Unit]
Description=GOrbitS API (Spring Boot)
After=network.target mariadb.service

[Service]
Type=simple
User=deploy
WorkingDirectory=/home/deploy/servers/gorbits
EnvironmentFile=/home/deploy/servers/gorbits/.env
ExecStart=/usr/bin/java -jar /home/deploy/servers/gorbits/current.jar --spring.profiles.active=mysql,prod
SuccessExitStatus=143
Restart=on-failure
RestartSec=10

[Install]
WantedBy=multi-user.target
